# services/gx-opendtu-zero-export/run
#!/bin/bash
exec 2>&1
exec python3 -u /data/gx-opendtu-zero-export/src/main.py --config /data/gx-opendtu-zero-export/config/config.json
