# services/gx-opendtu-zero-export/run
#!/bin/bash
exec 2>&1
cd /data/gx-opendtu-zero-export
exec python3 -u -m src.main --config /data/gx-opendtu-zero-export/config/config.json
